# server-files/home/apt-cacher/.config/systemd/user/update-registry.service
[Unit]
Description=Rebuild desktop image

[Service]
ExecStart=/home/apt-cacher/timer-scripts/update-registry.sh
