# server-files/home/nextcloud/.config/systemd/user/update-registry.service
[Unit]
Description=Rebuild desktop image

[Service]
ExecStart=/home/nextcloud/timer-scripts/update-registry.sh
