# server-files/etc/systemd/system/backup.service
[Unit]
Description=Backup /home

[Service]
ExecStart=/root/bin/backup.sh

